# launchd/com.paperflow.daemon.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <!-- Reference template. Prefer `paperflow install-daemon`, which fills the
         placeholders with your actual python path and repo location. -->
    <key>Label</key>
    <string>com.paperflow.daemon</string>
    <key>ProgramArguments</key>
    <array>
        <string>__PYTHON__</string>
        <string>-m</string>
        <string>paperflow.cli</string>
        <string>daemon</string>
    </array>
    <key>EnvironmentVariables</key>
    <dict>
        <key>PYTHONPATH</key>
        <string>__REPO__</string>
    </dict>
    <key>RunAtLoad</key>
    <true/>
    <key>KeepAlive</key>
    <true/>
    <key>StandardOutPath</key>
    <string>__HOME__/.paperflow/launchd.out.log</string>
    <key>StandardErrorPath</key>
    <string>__HOME__/.paperflow/launchd.err.log</string>
</dict>
</plist>
